# systemd/lyric-ui.service
[Unit]
Description=Lyric Speaker UI Service
After=lyric-ble.service network.target graphical.target pulseaudio.service
Wants=lyric-ble.service pulseaudio.service

[Service]
Type=simple
User=__USER__
Group=pulse-access
ExecStart=/opt/lyric-app/lyric_app ui
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lyric-ui

# 环境变量
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:/run/pulse/native
Environment=XDG_RUNTIME_DIR=/run/user/__UID__
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__UID__/bus

[Install]
WantedBy=multi-user.target
